# etc.Linux.MacBookPro11,1/systemd/system/fix-spurious-wakeups.service
[Service]
Type=oneshot
ExecStart=sh -c "echo XHC1 > /proc/acpi/wakeup"
RemainAfterExit=true
[Install]
WantedBy=multi-user.target
